# supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=python manage.py runserver 0.0.0.0:8000
directory=/app

[program:celery]
command=celery -A backend worker --loglevel=info
directory=/app

[program:bot]
command=python bot/main.py
directory=/app
